# curator4-ArtemI.conf
[program:curator4-ArtemI]
command=python3 /home/curator4/curator4-ArtemI/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-ArtemI/test.err.log
stdout_logfile=/home/curator4/curator4-ArtemI/test.out.log